# ansible/roles/inference_container/templates/rocm-inference-stack.service.j2
# Copyright (c) Advanced Micro Devices, Inc. All rights reserved.
#
# SPDX-License-Identifier: MIT

[Unit]
Description=ROCm inference stack (vLLM + LMCache + NIXL) container
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
Restart=on-failure
RestartSec=30
ExecStartPre=-/usr/bin/docker rm -f {{ inference_container_service_name }}
ExecStart=/usr/bin/docker run --name {{ inference_container_service_name }} \
  --device /dev/kfd --device /dev/dri \
  --network host --ipc host \
  --shm-size=10G \
  --ulimit memlock=-1 --ulimit stack=67108864 \
  {{ inference_container_image }} \
  sleep infinity
ExecStop=/usr/bin/docker stop -t 120 {{ inference_container_service_name }}

[Install]
WantedBy=multi-user.target
