# aai-day-release/ansible/roles/aai_day_stack/templates/aai-day-stack.service.j2
[Unit]
Description=aai-day inference stack (vLLM + LMCache)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
Restart=on-failure
RestartSec=30
WorkingDirectory={{ aai_day_compose_dir }}
ExecStart=/usr/bin/docker compose --env-file {{ aai_day_compose_dir }}/.env up
ExecStop=/usr/bin/docker compose --env-file {{ aai_day_compose_dir }}/.env down

[Install]
WantedBy=multi-user.target
